# files/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null

[include]
files = /etc/supervisor/conf.d/*.conf

[program:simplemonitor]
command=python %(ENV_MONITOR_FOLDER)s/monitor.py -f %(ENV_MONITOR_CONFIG_FILE)s %(ENV_MONITOR_ARGS)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=1
startsecs=3
autorestart=true
stopsignal=HUP

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=5

